# systemd/rapiba-backup.service
[Unit]
Description=Rapiba Scheduled Backup
Documentation=https://github.com/yourusername/rapiba
PartOf=rapiba-backup.timer

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/rapiba --backup-now
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rapiba-scheduled
