# systemd/rapiba-backup.timer
[Unit]
Description=Rapiba Scheduled Backup Timer
Documentation=https://github.com/yourusername/rapiba

[Timer]
# Täglich um 2:00 Uhr nachts
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# Oder stündlich
# OnCalendar=hourly

# Oder alle 6 Stunden
# OnCalendar=*-*-* 00,06,12,18:00:00

# Bei Boot starten wenn Zeitpunkt verpasst wurde
Persistent=true

[Install]
WantedBy=timers.target
